# deploy/systemd/micro-niche-auto-seeds.service
[Unit]
Description=Micro Niche Finder auto seed and report generator
After=network.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/micro-niche-finder
EnvironmentFile=/etc/micro-niche-finder.env
ExecStart=/opt/micro-niche-finder/.venv/bin/python -m apps.worker.bootstrap_auto_seeds --seed-count 5 --candidate-count 5 --top-k 3
