# deploy/systemd/micro-niche-daily-report.timer
[Unit]
Description=Run Micro Niche Finder daily report delivery at 09:00 KST

[Timer]
OnCalendar=*-*-* 09:00:00 Asia/Seoul
Unit=micro-niche-daily-report.service
Persistent=true

[Install]
WantedBy=timers.target
