# deploy/systemd/micro-niche-kosis-collector.service
[Unit]
Description=Micro Niche Finder KOSIS employee-count collector
After=network.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/micro-niche-finder
EnvironmentFile=/etc/micro-niche-finder.env
ExecStart=/opt/micro-niche-finder/.venv/bin/python -m apps.worker.run_kosis_collector
