# deploy/systemd/micro-niche-naver-shopping-insight-collector.timer
[Unit]
Description=Run Micro Niche Finder Naver Shopping Insight collector once per day

[Timer]
OnCalendar=*-*-* 01:40:00 Asia/Seoul
AccuracySec=1m
Unit=micro-niche-naver-shopping-insight-collector.service
Persistent=true

[Install]
WantedBy=timers.target
