# deploy/systemd/micro-niche-collector.timer
[Unit]
Description=Run Micro Niche Finder collector every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Unit=micro-niche-collector.service
Persistent=true

[Install]
WantedBy=timers.target
