# deploy/systemd/micro-niche-daily-report.service
[Unit]
Description=Micro Niche Finder daily report delivery
After=network.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/micro-niche-finder
EnvironmentFile=/etc/micro-niche-finder.env
ExecStart=/opt/micro-niche-finder/.venv/bin/python -m apps.worker.run_daily_report --refresh-seeds
