# deploy/systemd/micro-niche-naver-search-collector.service
[Unit]
Description=Micro Niche Finder Naver Search collector
After=network.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/micro-niche-finder
EnvironmentFile=/etc/micro-niche-finder.env
ExecStart=/usr/bin/flock /tmp/micro-niche-finder-db.lock /opt/micro-niche-finder/.venv/bin/python -m apps.worker.run_naver_search_collector
